# s6/services/traefik.svc/run
#!/bin/sh
exec 2>&1
cd /bin
echo -n "Starting Traefik"
exec ./traefik
